# rpms/common/environment/awips2-cave/profile.d/awips2Cave.csh
#!/bin/csh

# Determine where cave has been installed.
set CAVE_INSTALL="/awips2/cave"
setenv TMCP_HOME "${CAVE_INSTALL}/caveEnvironment"
setenv FXA_HOME "${CAVE_INSTALL}/caveEnvironment"
