# rpms/awips2.core/Installer.httpd-pypies/configuration/usr/lib/systemd/system/httpd-pypies.service.d/httpd-pypies.conf
[Service]
Environment="FXA_HOME=/awips2/cave/caveEnvironment"
Environment="YAJSW_HOME=/awips2/yajsw"
Environment="LD_LIBRARY_PATH=/awips2/ldm/lib:/awips2/tools/lib:/awips2/python/lib:"
Environment="PATH=/awips2/ldm/bin:/awips2/ldm/decoders:/awips2/ldm/util:/awips2/java/bin:/awips2/GFESuite/bin:/awips2/python/bin:/awips2/groovy/bin:/awips2/tools/bin:/usr/local/bin:/bin:/usr/bin:/usr/local/sbin:/usr/sbin:/home/awips/.local/bin:/home/awips/bin"
Environment="JAVA_HOME=/awips2/java"
Environment="TMCP_HOME=/awips2/cave/caveEnvironment"
Environment="GROOVY_HOME=/awips2/groovy"
Environment="LDMHOME=/awips2/ldm"
