# rpms/awips2.core/Installer.edex-environment/wes2bridge.files/deliverables/profile.d/awips2-edex-env.csh
#!/bin/csh

set EDEX_ENV_PATH="/awips2/edex-environment/macro"

if $?PATH then
   setenv PATH ${EDEX_ENV_PATH}:$PATH
else
   setenv PATH ${EDEX_ENV_PATH}
endif
