# rpms/awips2.core/Installer.rcm/scripts/init.d/edex_rcm
#!/bin/bash
#
# edex_rcm  This scripts takes care of starting and stopping the
#           radar comms manager.
# chkconfig: - 99 10
#
# description: Radar comms manager control scripts

# Source function library.
. /etc/rc.d/init.d/functions

# Source networking configuration.
. /etc/sysconfig/network

# Check that networking is up.
[ ${NETWORKING} = "no" ] && exit 0

RETVAL=0

# user to run the Radar Server
EDEXUSER=awips

# the java installation.
JAVA_INSTALL="/awips2/java"
# the Radar Server installation.
RCM_INSTALL="/awips2/rcm"

# set java location
JAVA_HOME_DIR=${JAVA_INSTALL}
export PATH=${JAVA_HOME_DIR}/bin:$PATH

# set Radar Server install directory
RCM_HOME=${RCM_INSTALL}
RCM_DATA=$RCM_HOME/data

start() {
        su ${EDEXUSER} -l -c "sh ${RCM_HOME}/bin/start" > /dev/null 2>&1
}

stop() {
        su ${EDEXUSER} -l -c "sh ${RCM_HOME}/bin/stop" > /dev/null 2>&1
}

# Verify root user
checkUser() {
   REQUIREDUSER="root"
   CURUSER=`whoami`
   if [ "$CURUSER" != "$REQUIREDUSER" ]; then
      echo "Insufficient privileges: must run script as $REQUIREDUSER"
      exit 1
   fi
}

# See how we were called.
case $1 in
  start)
    checkUser
        echo -n "Starting RCM: "
        start
        RETVAL=$?
        if [ -e ${RCM_DATA}/radarserver-pid ]; then
            echo "OK"
        else
            echo "FAILED"
        fi
        ;;
  stop)
    checkUser
        echo -n "Stopping RCM: "
        stop
        RETVAL=$?
        if [ -e ${RCM_DATA}/radarserver-pid ]; then
            echo "FAILED"
        else
            echo "OK"
        fi
        ;;
  status)
        status -p ${RCM_DATA}/radarserver-pid RadarServer
        RETVAL=$?
        ;;
  restart)
    checkUser
	echo -n "Restarting RCM: "
	stop
	start
	RETVAL=$?
	;;
  *)
	# Print help
	echo "Usage: $0 {start|stop|status|restart}" 1>&2
	exit 1
	;;
esac

exit $RETVAL
