# rpms/awips2.core/Installer.groovy/scripts/profile.d/awips2Groovy.sh
#!/bin/bash

export GROOVY_HOME=/awips2/groovy
export PATH=${GROOVY_HOME}/bin:${PATH}
