# rpms/awips2.core/Installer.yajsw/scripts/profile.d/awips2Yajsw.sh
#!/bin/bash

if [ -d /awips2/yajsw ]; then
   YAJSW_INSTALL="/awips2/yajsw"
# Update The Environment
   export YAJSW_HOME=${YAJSW_INSTALL}
fi
